# debian/prerm
#!/bin/sh
case "$1" in
remove|purge)
	invoke-rc.d xrdp stop
	update-rc.d xrdp remove
	;;
esac
exit 0
